# php-server.conf
[program:php-server]tes=1MB
user=rootgfile_backups=0
directory=/root/code/php-server
command=/usr/bin/php -S 0.0.0.0:3218
process_name=%(program_name)s
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/kcptun/server.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=0