# shadowsocks.conf
[program:shadowsocks]
user=root
directory=/usr/local/bin
command=ssserver -p 500 -k chan -m aes-256-cfb
process_name=%(program_name)s
autostart=true